# rootfs/etc/services.d/syslog/run
#!/usr/bin/with-contenv sh
exec syslogd -n -D -O /dev/stdout
